# autodeploy.service
[Unit]
Description=Auto deploy from GitHub

[Service]
Type=oneshot
ExecStart=/root/autodeploy.sh
